# modules/helloworld/files/helloworld.conf
description "helloworld.js server"
author      "agraham - betable"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    # Not sure why $HOME is needed, but we found that it is:
    export HOME="/root"

    exec /usr/bin/node /helloworld/helloworld.js >> /var/log/node.log 2>&1
end script
